# install/ricettario.service
[Unit]
Description=A self-hosted web database of family recipes

[Service]
ExecStart=/usr/local/bin/ricettario

[Install]
WantedBy=multi-user.target